# crontab
10 * * * * ~/photoScript.sh >> ~/photolog.txt
40 * * * * ~/photoScript.sh >> ~/photolog.txt

